# host-gate/systemd/host-gate.service
[Unit]
Description=Host Gate - Container Command Approval Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/host-gate-daemon
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t/host-gate %h/.local/share/ykman
PrivateTmp=true

[Install]
WantedBy=default.target
